# supervisord.conf
[program:server]
command=node ./bin/server.js
