# supervisord.conf
[program:daisy]
command=node ./bin/server.js
